# sys-apps/usbguard/files/usbguard.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "USBGuard - a service that handles enforcing USB device policy."
author          "allenwebb@chromium.org"

start on screen-locked
stop on screen-unlocked or start-user-session

respawn
respawn limit 15 5

umask 066

exec /sbin/minijail0 -S /opt/google/usbguard/usbguard-daemon-seccomp.policy \
    /usr/sbin/usbguard-daemon -C -P -s
