# net-misc/rmtfs/files/77-rmtfs.rules
ACTION=="add", KERNEL=="qcom_rmtfs_mem?", GROUP=="rmtfs", MODE="0660", RUN+="/sbin/initctl emit --no-wait qcom-rmtfs-added"