# chromeos-base/swap-init/files/init/swap.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Setup compressed memory swap"
author          "chromium-os-dev@chromium.org"

env MIN_LOW_MEMORY_MARGIN=

start on starting boot-services
task

exec /usr/share/cros/init/swap.sh start
