# chromeos-base/swap-init/files/init/teardown-swap.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Teardown swap on shutdown"
author          "chromium-os-dev@chromium.org"

start on stopping boot-services
task

exec /usr/share/cros/init/swap.sh stop
