# net-print/cups/files/init/cups-clear-state.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Perform CUPS session cleanup tasks."
author          "chromium-os-dev@chromium.org"

start on login-prompt-visible
task

script
  if grep -q CHROMEOS_RELEASE_TRACK=testimage /etc/lsb-release && [ -f /run/cups/disable ]; then
    logger -t "${UPSTART_JOB}" "cups-clear-state disabled for testing."
  else
    exec /usr/share/cros/init/cups-clear-state.sh
  fi
end script
