# net-dns/avahi-daemon/files/init/cgroup.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Manages CPU cgroup for avahi daemon"
author          "chromium-os-dev@chromium.org"

start on started avahi
stop on stopped avahi

oom score 0

env CGROUP_DIR=/sys/fs/cgroup/cpu/avahi-daemon

script
  mkdir -p "${CGROUP_DIR}"  # croslint: disable
  echo 1000000 > "${CGROUP_DIR}/cpu.cfs_period_us"
  echo 5000 > "${CGROUP_DIR}/cpu.cfs_quota_us"
  for p in $(pgrep avahi); do echo "$p" > "${CGROUP_DIR}/cgroup.procs"; done
end script
